# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:relay]
command=node /app/relay/server.js
directory=/app/relay
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=OUTPUT="%(ENV_OUTPUT)s",PROFILE="%(ENV_PROFILE)s",WIDTH="%(ENV_WIDTH)s",HEIGHT="%(ENV_HEIGHT)s",FRAMERATE="%(ENV_FRAMERATE)s",WS_PORT="%(ENV_WS_PORT)s",CHANNEL_NAME="%(ENV_CHANNEL_NAME)s",CHANNEL_ID="%(ENV_CHANNEL_ID)s",PROGRAMME_TITLE="%(ENV_PROGRAMME_TITLE)s",PROGRAMME_DESC="%(ENV_PROGRAMME_DESC)s",STREAM_URL="%(ENV_STREAM_URL)s"
priority=10

[program:chrome]
command=bash /tmp/launch-chrome.sh
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=URL="%(ENV_URL)s",WIDTH="%(ENV_WIDTH)s",HEIGHT="%(ENV_HEIGHT)s",CDP_PORT="%(ENV_CDP_PORT)s",EXTENSION_ID="%(ENV_EXTENSION_ID)s",EXTENSION_DIR="%(ENV_EXTENSION_DIR)s"
priority=20

[program:trigger]
command=bash /app/trigger-capture.sh
autostart=true
autorestart=unexpected
startretries=5
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CDP_PORT="%(ENV_CDP_PORT)s",WS_PORT="%(ENV_WS_PORT)s",WIDTH="%(ENV_WIDTH)s",HEIGHT="%(ENV_HEIGHT)s",FRAMERATE="%(ENV_FRAMERATE)s"
priority=30
